# launchd/com.thinkwork.factory.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  TEMPLATE — do NOT load this file as-is. `factoryd install` renders every
  double-underscore placeholder to an ABSOLUTE path and writes the result into
  the user's Library/LaunchAgents directory, then bootstraps it with launchctl
  in the per-user GUI domain.

  launchd NEVER sources a shell rc, so it has no PATH and no notion of `node`:
  a bare `node` or a home-relative path here is the #1 LaunchAgent failure.
  Both slots below are substituted with absolute paths, PATH is set explicitly.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>__LABEL__</string>

  <key>ProgramArguments</key>
  <array>
    <string>__NODE_BIN__</string>
    <string>__ENTRY__</string>
    <string>run</string>
  </array>

  <key>WorkingDirectory</key>
  <string>__WORKING_DIR__</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>__PATH__</string>
    <key>THINKWORK_FACTORY_DIR</key>
    <string>__STATE_DIR__</string>
    <key>NODE_ENV</key>
    <string>production</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
  </dict>

  <key>ThrottleInterval</key>
  <integer>__THROTTLE_INTERVAL__</integer>

  <key>ProcessType</key>
  <string>Background</string>

  <key>StandardOutPath</key>
  <string>__LOG_PATH__</string>
  <key>StandardErrorPath</key>
  <string>__LOG_PATH__</string>
</dict>
</plist>
